# scripts/crontab
* * * * * echo "{}" | kcat -P -b kafka-1:9094 -t check